# deploy/systemd/kronoskvm-containers.service
[Unit]
Description=KronosKVM application containers
Requires=docker.service
After=docker.service network-online.target
Wants=network-online.target
Conflicts=kronoskvm-api.service

[Service]
Type=oneshot
RemainAfterExit=yes
WorkingDirectory=/opt/kronoskvm
ExecStartPre=/usr/bin/install -d -m 0755 -o root -g root /mnt/kronoskvm-storage
ExecStart=/opt/kronoskvm/scripts/start-containers.sh
ExecReload=/opt/kronoskvm/scripts/start-containers.sh
ExecStop=/opt/kronoskvm/scripts/stop-containers.sh
TimeoutStartSec=600
TimeoutStopSec=120

[Install]
WantedBy=multi-user.target
