# systemd/kronoskvm.service
[Unit]
Description=KronosKVM control service
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=kronoskvm
Group=kronoskvm
ExecStart=/opt/kronoskvm/.venv/bin/kronoskvm --config /etc/kronoskvm/config.toml
Restart=on-failure
RestartSec=3

NoNewPrivileges=true
PrivateTmp=true
ProtectHome=true
ProtectSystem=strict
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true

[Install]
WantedBy=multi-user.target
